# conf-files/supervisor.conf
[program:<name>]
environment = <environment>
command = uwsgi --ini <uwsgi_ini_file>
autostart=true
